# supervisord.conf
[supervisord]
nodaemon=true

[program:worker]
user=root
command=ruby -I /app/lib/ /app/bin/start
directory=/app
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/warehouse/warehouse.log
